# systemd/gemini-telegram-bot.service
[Unit]
Description=Gemini Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wilsonwkj
WorkingDirectory=/home/wilsonwkj/Workspace/gemini-telegram-bot
EnvironmentFile=/home/wilsonwkj/Documents/secrets/gemini-telegram-agent.env
Environment=PATH=/home/wilsonwkj/.local/bin:/home/wilsonwkj/.npm-global/bin:/home/wilsonwkj/.nvm/versions/node/v22.22.2/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
ExecStart=/home/wilsonwkj/Workspace/gemini-telegram-bot/.venv/bin/python /home/wilsonwkj/Workspace/gemini-telegram-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gemini-telegram-bot

NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=default.target
